# deployment/romm2steamdeck.service
[Unit]
Description=RomM2SteamDeck - Download ROMs from RomM to Steam Deck
After=network.target

[Service]
Type=simple
# For user services, paths are relative to the user's home directory
# Adjust these paths to match your installation:

# For source installation with Python:
ExecStart=/usr/bin/python3 /path/to/app.py
# Or for AppImage:
# ExecStart=/home/deck/RomM2SteamDeck-x86_64.AppImage
# Or for virtual environment:
# ExecStart=/path/to/venv/bin/python /path/to/app.py

# Working directory (adjust path as needed)
WorkingDirectory=%h/.config/romm2steamdeck

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=romm2steamdeck

# Security (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
